# [project start]/jsub.service
[Unit]
Description=The Msg Soki Service

[Service]
Restart=on-failure
RestartSec=5s
ExecStart=node /var/www/sub/back.index.js

#3359